# deploy/systemd/moex-scalper-dashboard.service
[Unit]
Description=MOEX scalper dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__RUN_USER__
Group=__RUN_GROUP__
WorkingDirectory=__PROJECT_DIR__
ExecStart=/usr/bin/env bash __PROJECT_DIR__/scripts/run_dashboard_service.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
