# deploy/systemd/moex-scalper-govern.timer
[Unit]
Description=Nightly coordinated governance for MOEX scalper

[Timer]
OnCalendar=Mon..Fri 18:24:00
Persistent=true
Unit=moex-scalper-govern.service

[Install]
WantedBy=timers.target
