# deploy/systemd/moex-scalper-summary.timer
[Unit]
Description=Nightly daily summary report for MOEX scalper

[Timer]
OnCalendar=Mon..Fri 18:26:00
Persistent=true
Unit=moex-scalper-summary.service

[Install]
WantedBy=timers.target
