# deploy/systemd/moex-scalper-tune.timer
[Unit]
Description=Nightly paper autotune for MOEX scalper

[Timer]
OnCalendar=Mon..Fri 18:24:00
Persistent=true
Unit=moex-scalper-tune.service

[Install]
WantedBy=timers.target
